# scripts/launchd/com.zfit.watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.zfit.watcher</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__HOME__/dev/zFit/scripts/watch-deploy.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__HOME__/dev/zFit</string>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/zfit-watcher.log</string>

  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/zfit-watcher.log</string>
</dict>
</plist>
